# templates/default/app.upstart.erb
description "Start the app server"

start on startup
stop on shutdown

respawn
respawn limit 10 5


setuid <%= @app_user %>

chdir <%= @working_dir %>

exec \
gunicorn <%= @app_id %>.wsgi:application -b :8000 >> <%= @working_dir %>/logs/gunicorn.sys.log 2>&1
